# kbsurfer.conf
# kbsurfer.conf

description     "Start kbsurfer.com"

# Start up when the system hits any normal runlevel, and
# shuts down when the system goes to shutdown or reboot.
start on filesystem or runlevel [2345]
stop on runlevel [06]

# IMPORTANT: You will want to use this with Forever. It
# tells Upstart that forever is going to fork after it
# starts.
expect fork

# This monitors Forever, which seems gratuitous.
# TIP: Comment these out while debugging your script.
respawn
respawn limit 5 30

# Send error messages to the console. Useful for debugging.
console output

# The user's home directory
env HOME=/var/kbsurfer.com

# Now we run Forever, telling it to write its logs to
# /tmp, and to run the script /opt/example/server.js
script
        cd $HOME
        exec forever start -a -l $HOME/logs/kbsurfer.com.log -o $HOME/logs/kbsurfer.com.stdout.log -e $HOME/logs/kbsurfer.com.stderr.log app.js -p 80
end script
